# crontab
5 * * * * cat /var/log/asterisk/messages | /usr/bin/php --process-file ~/.local/almon/almon.php >> /tmp/asterisk.log
